# service/xserver.service
[Unit]
Description=naked X server
After=ssh.server

[Service]
Type=simple
ExecStartPre=export DISPLAY=:0; xset s off
ExecStart=/usr/bin/X -nolock -nocursor :0
KillMode=process
Restart=on-failure
RestartPreventExitStatus=255

[Install]
WantedBy=multi-user.target
Alias=sshd.service
